# config/systemd/headless-siem-ruled.service
[Unit]
Description=Headless SIEM - Sigma Rule Engine
Documentation=https://github.com/your-org/headless-siem
After=network.target headless-siem-pipes.service headless-siem-normalized.service
Requires=headless-siem-pipes.service
Wants=headless-siem-normalized.service

[Service]
Type=simple
User=user
Environment=HEADLESS_SIEM_ROOT=/home/user/projects/headless-siem
WorkingDirectory=/home/user/projects/headless-siem
# Suppress known false positives per config/rules/suppress.toml (see that
# file's header for the rule grammar / opt-in note).
ExecStart=/home/user/projects/headless-siem/target/release/ruled \
  --rules /home/user/projects/headless-siem/config/rules \
  --output /home/user/projects/headless-siem/data/alerts \
  --suppress /home/user/projects/headless-siem/config/rules/suppress.toml
StandardInput=file:/run/headless-siem/normalized-out.pipe
StandardOutput=file:/run/headless-siem/ruled-out.pipe
StandardError=journal
# StandardInput is a pipe fed by normalized: any normalized restart closes
# the write end, ruled hits EOF and exits 0 (clean, not a "failure"), so
# Restart=on-failure silently never brings it back. Restart=always covers
# that case too (see implementation-plan.md issue #46).
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/user/projects/headless-siem/data/alerts /run/headless-siem
ReadOnlyPaths=/home/user/projects/headless-siem/config

[Install]
WantedBy=multi-user.target
