# config/systemd/headless-siem-retention.timer
[Unit]
Description=Headless SIEM - Nightly Retention Cleanup Timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
