# config/systemd/headless-siem-normalized.service
[Unit]
Description=Headless SIEM - Log Normalizer
Documentation=https://github.com/your-org/headless-siem
After=network.target headless-siem-pipes.service
Requires=headless-siem-pipes.service

# normalized listens directly on UDP/TCP (config/normalized.toml's [listen])
# and writes normalized JSONL to a named pipe consumed by ruled. No local
# rsyslog relay — clawserver doesn't have rsyslog installed, and this
# standalone-listener mode is the same one already proven reliable for
# several days of real traffic via dev.sh, just running under systemd
# instead. Senders (pfSense, Proxmox, jumphost) forward straight here; see
# llm-based-soc/documentation/log-forwarding-setup.md.
#
# Alternative — local rsyslog relay via a named pipe (disk-backed queue,
# useful if you want rsyslog's own buffering/durability in front of
# normalized): see config/rsyslog.d/50-headless-siem.conf, then
#   1. Add --stdin to ExecStart below
#   2. Add StandardInput=file:/run/headless-siem/normalized-in.pipe
#   3. Copy 50-headless-siem.conf to /etc/rsyslog.d/ and restart rsyslog

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/projects/headless-siem

ExecStart=/home/user/projects/headless-siem/target/release/normalized \
  --data-dir /home/user/projects/headless-siem/data \
  --config /home/user/projects/headless-siem/config/normalized.toml

StandardOutput=file:/run/headless-siem/normalized-out.pipe
StandardError=journal

# Observed 2026-07-17: normalized clean-exited (exit 0, no crash/panic/OOM
# in the journal) after ~81min uptime with Restart=on-failure silently
# *not* restarting it, since a 0 exit isn't a "failure" — the whole
# ingestion pipeline sat blind until manually discovered. Root cause of
# the exit itself is still unknown, but per the same reasoning already
# applied to ruled/correlated (issue #46), a log normalizer has no
# legitimate reason to exit on its own, so Restart=always at least turns
# a silent death into a fast, visible restart-loop instead of permanent
# silence.
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/user/projects/headless-siem/data /run/headless-siem
ReadOnlyPaths=/home/user/projects/headless-siem/config

[Install]
WantedBy=multi-user.target
